# infra/systemd/tg-token-verifier.service
[Unit]
Description=tg-token v2 — EIP-1271/6492 signature verifier (viem)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/claude_folder/claude/tg-token/webapp_verifier
EnvironmentFile=-%h/Desktop/claude_folder/claude/tg-token/.env
Environment=PORT=8090 HOST=127.0.0.1
ExecStart=/usr/bin/node node_modules/.bin/tsx src/server.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
