# infra/systemd/tg-token-api.service
[Unit]
Description=tg-token v2 — FastAPI server (Mini App backend + webhook receivers)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/claude_folder/claude/tg-token
EnvironmentFile=%h/Desktop/claude_folder/claude/tg-token/.env
ExecStart=%h/Desktop/claude_folder/claude/tg-token/.venv/bin/uvicorn app.api:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
